# deploy/postgres/bw-monitor-maintenance-watchdog.timer
[Unit]
Description=VirtInfra Monitor maintenance queue watchdog

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=bw-monitor-maintenance-dispatch.service

[Install]
WantedBy=timers.target
