# deploy/postgres/bw-monitor-maintenance@.service
[Unit]
Description=VirtInfra Monitor PostgreSQL maintenance job %i
After=network-online.target docker.service bw-monitor.service
Wants=network-online.target docker.service
ConditionPathExists=/opt/bw-monitor/maintenance.py

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/bw-monitor
EnvironmentFile=/etc/default/bw-monitor
Environment=PYTHONUNBUFFERED=1
Environment=BW_MONITOR_APP_FILE=/opt/bw-monitor/app.py
Environment=BW_MONITOR_APP_SERVICE=bw-monitor.service
ExecStart=/opt/bw-monitor/venv/bin/python3 /opt/bw-monitor/maintenance.py %i
ExecStopPost=/bin/bash -lc 'marker=/run/bw-monitor-maintenance-web-offline; if [ -s "$marker" ]; then state=$(systemctl is-system-running 2>/dev/null || true); svc=$(head -n1 "$marker" | tr -cd "A-Za-z0-9@_.:-"); case "$state" in stopping|offline) ;; *) [ -n "$svc" ] && systemctl reset-failed "$svc" >/dev/null 2>&1 || true; [ -n "$svc" ] && systemctl start "$svc" >/dev/null 2>&1 || true ;; esac; rm -f "$marker"; fi'
TimeoutStartSec=12h
TimeoutStopSec=60
KillSignal=SIGTERM
SendSIGKILL=yes
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=journal
StandardError=journal
