# deploy/postgres/bw-monitor.service
[Unit]
StartLimitIntervalSec=0
Description=VirtInfra Monitor PostgreSQL/TimescaleDB web and ingestion service
Wants=network-online.target docker.service
After=network-online.target docker.service
Requires=docker.service
ConditionPathExists=/opt/bw-monitor/app.py
ConditionPathExists=/opt/bw-monitor/venv/bin/gunicorn

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bw-monitor
EnvironmentFile=/etc/default/bw-monitor
Environment=PYTHONUNBUFFERED=1
RuntimeDirectory=bw-monitor
RuntimeDirectoryMode=0750
UMask=0077
ExecStartPre=/opt/bw-monitor/venv/bin/python3 -m py_compile /opt/bw-monitor/app.py /opt/bw-monitor/bw_pg.py /opt/bw-monitor/maintenance_native.py
ExecStartPre=/opt/bw-monitor/venv/bin/python3 -c "import bw_pg; print(bw_pg.healthcheck())"
ExecStart=/opt/bw-monitor/start-monitor.sh
Restart=always
RestartSec=2
TimeoutStartSec=180
TimeoutStopSec=60
KillSignal=SIGTERM
LimitNOFILE=1048576
TasksMax=32768
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/bw-monitor /var/lib/bw-monitor /run
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
Alias=virtinfra-monitor.service
