# deploy/postgres/bw-monitor-retention.service
[Unit]
Description=VirtInfra Monitor exact 2-day raw / 7-day hourly PostgreSQL retention
After=network-online.target docker.service bw-monitor.service
Wants=network-online.target docker.service
ConditionPathExists=/opt/bw-monitor/app.py
ConditionPathExists=/opt/bw-monitor/retention.py
ConditionPathExists=/opt/bw-monitor/maintenance_native.py

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/bw-monitor
EnvironmentFile=/etc/default/bw-monitor
Environment=PYTHONUNBUFFERED=1
Environment=BW_MONITOR_APP_FILE=/opt/bw-monitor/app.py
ExecStart=/opt/bw-monitor/venv/bin/python3 /opt/bw-monitor/retention.py
TimeoutStartSec=12h
TimeoutStopSec=60
KillSignal=SIGTERM
SendSIGKILL=yes
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=journal
StandardError=journal
